# xci/source/info.apps/Manage-Outages/tags/tag-1.7-20191111/sbin/V4Outage.crontab
1,16,31,46 * * * * /bin/bash -l -c "/soft/warehouse-apps-1.0/Manage-Outages/PROD/sbin/V4Outage.cron"
